# debian/abraflexi-server-backup.cron.daily
#!/bin/bash
. /etc/default/flexibee-server-backup
/usr/share/flexibee/bin/flexibee-backup ${BACKUPDIR}

