# debian/flexibee-server.conffiles
/etc/logrotate.d/flexibee
/etc/ufw/applications.d/flexibee
/etc/default/flexibee
/etc/init.d/flexibee
/etc/flexibee/flexibee-server.xml
